# scripts/launchd/com.rileyfile.contextsync.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.rileyfile.contextsync</string>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>/Users/rileycolleyFW/dev/rileyfile/scripts/sync_runner.sh</string>
  </array>

  <key>WorkingDirectory</key>
  <string>/Users/rileycolleyFW/dev/rileyfile</string>

  <key>EnvironmentVariables</key>
  <dict>
    <key>TZ</key>
    <string>America/New_York</string>
  </dict>

  <key>StartCalendarInterval</key>
  <dict>
    <key>Hour</key>
    <integer>23</integer>
    <key>Minute</key>
    <integer>30</integer>
  </dict>

  <key>StartInterval</key>
  <integer>300</integer>

  <key>WatchPaths</key>
  <array>
    <string>/Users/rileycolleyFW/.codex</string>
    <string>/Users/rileycolleyFW/.claude</string>
    <string>/Users/rileycolleyFW/.gemini</string>
    <string>/Users/rileycolleyFW/dev</string>
    <string>/Users/rileycolleyFW/Documents/GitHub</string>
    <string>/Users/rileycolleyFW/Desktop/HarmonyApp</string>
    <string>/Users/rileycolleyFW/Library/Mobile Documents/com~apple~CloudDocs/RileyFile</string>
    <string>/Users/rileycolleyFW/Library/Group Containers/group.com.apple.notes/NoteStore.sqlite</string>
    <string>/Users/rileycolleyFW/Library/Group Containers/9K33E3U3T4.net.shinyfrog.bear/Application Data/database.sqlite</string>
    <string>/Users/rileycolleyFW/dev/rileyfile/inbox/chatgpt</string>
    <string>/Users/rileycolleyFW/dev/rileyfile/inbox/gemini</string>
    <string>/Users/rileycolleyFW/dev/rileyfile/inbox/claude</string>
    <string>/Users/rileycolleyFW/dev/rileyfile/inbox/codex</string>
  </array>

  <key>StandardOutPath</key>
  <string>/Users/rileycolleyFW/dev/rileyfile/logs/contextsync.stdout.log</string>

  <key>StandardErrorPath</key>
  <string>/Users/rileycolleyFW/dev/rileyfile/logs/contextsync.stderr.log</string>

  <key>RunAtLoad</key>
  <true/>
</dict>
</plist>
